# network/if-up.d/static-arp
#!/bin/sh
# Script to add static ARP entry for default gateway
sudo /usr/sbin/arp -s 193.191.176.254 ca:fe:c0:ff:ee:00 -i eth0
